# scripts/month.every1.cron
* * 1 * * /home/xilibo/php/scripts/cron_minute.php
